# launchd/com.thepregames.kings.frontend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.thepregames.kings.frontend</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/thurston/Workspace/Kings/scripts/run_frontend.sh</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>WorkingDirectory</key>
  <string>/Users/thurston/Workspace/Kings/dist</string>
  <key>StandardOutPath</key>
  <string>/Users/thurston/Workspace/Kings/logs/frontend.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/thurston/Workspace/Kings/logs/frontend.log</string>
</dict>
</plist>
